# docker/prod/supervisor/queue.conf

[program:queue-work]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/html/artisan queue:work
autostart=true
autorestart=true
numprocs=5
